# opencanary-installer/opencanary-heartbeat.timer
[Unit]
Description=opencanary heartbeat

[Timer]
OnCalendar=*-*-* *:*:00
Unit=opencanary-heartbeat.service

[Install]
WantedBy=timers.target
